# device/systemd/marvin-temp-guard.service
[Unit]
Description=Marvin device SoC 溫度防護（過熱→本機靜音 DigiAMP）
After=network-online.target
[Service]
Environment=MARVIN_TEMP_STOP_C=50
Environment=MARVIN_TEMP_CLEAR_C=48
ExecStart=/usr/local/bin/temp_guard.sh
Restart=always
RestartSec=3
[Install]
WantedBy=multi-user.target
