# device/systemd/marvin-temp-guard.service
[Unit]
Description=Marvin device SoC жә«еәҰйҳІиӯ·пјҲйҒҺзҶұвҶ’йҖҒеҒңжӯўж’ӯж”ҫзөҰеӨ§и…Ұпјү
After=network-online.target
[Service]
Environment=MARVIN_TEXT_TOKEN=__MARVIN_TEXT_TOKEN__
Environment=MARVIN_MAC_SAY_URL=http://100.123.68.86:8790/say
Environment=MARVIN_TEMP_STOP_C=65
Environment=MARVIN_TEMP_CLEAR_C=57
ExecStart=/usr/local/bin/temp_guard.sh
Restart=always
RestartSec=3
[Install]
WantedBy=multi-user.target
