# device/systemd/wyoming-satellite.service.d/override.conf
# 場景＝外出：XVF3800 當麥克風 + BT(Soundcore) 當喇叭。單一輸出、不做雙路 AEC。
#   mic  = arecord XVF3800(Array)；XVF3800 未接時 arecord crash-loop，插回會自動重連恢復。
#   snd  = aplay -D btspk（asound.conf 的 bluealsa a2dp Soundcore）；Soundcore 未連時無聲。
# 家用 AirPlay 走另一條（shairport → DigiAMP），與此無關。
# override 整段蓋掉 base unit，故 satellite.env 的 MIC_DEVICE/SND_DEVICE 不被讀（換裝置直接改本檔）。
[Unit]
After=
After=network-online.target
Wants=
Requires=

[Service]
Environment=HOME=/home/pi
ExecStartPre=
ExecStart=
ExecStart=/home/pi/wyoming-satellite/.venv/bin/python -u -m wyoming_satellite --name marvin-satellite --uri tcp://0.0.0.0:10700 --mic-command "arecord -D plughw:CARD=Array,DEV=0 -r 16000 -c 1 -f S16_LE -t raw" --mic-command-rate 16000 --mic-command-channels 1 --mic-command-width 2 --snd-command "aplay -D btspk -r 48000 -c 2 -f S16_LE -t raw" --snd-command-rate 48000 --snd-command-channels 2 --snd-command-width 2
