# device/systemd/wyoming-satellite.service.d/override.conf
# 麥克風已移除 → 純輸出音箱。用節流靜音源取代失敗的 arecord（避免 crash-loop 洗 log），
# 保留麥克風幀流讓 Mac 大腦 Sentinel 心跳持續看到活動。不再依賴 openwakeword。
[Unit]
After=
After=network-online.target
Wants=
Requires=

[Service]
ExecStartPre=
ExecStart=
ExecStart=/home/pi/wyoming-satellite/.venv/bin/python -u -m wyoming_satellite --name marvin-satellite --uri tcp://0.0.0.0:10700 --mic-command "sh -c 'while :; do head -c 3200 /dev/zero; sleep 0.1; done'" --mic-command-rate 16000 --mic-command-channels 1 --mic-command-width 2 --snd-command "aplay -D ${SND_DEVICE} -r ${SND_RATE} -c ${SND_CHANNELS} -f S16_LE -t raw" --snd-command-rate ${SND_RATE} --snd-command-channels ${SND_CHANNELS} --snd-command-width 2
