# device/systemd/marvin-temp-guard.service.d/override.conf
# 排線把 Pi 與 DigiAMP 分板後，SoC 與 AMP 熱脫鉤：SoC 靠自身閒置熱在盒內
# 爬到 52-58°C，原本 STOP=50 會播到一半誤靜音，且此門檻已無法反映 AMP 溫度。
# 拉高到 60 消除誤觸發（SoC 54°C idle 離降頻 80°C 尚有 26°C 裕度）。
# 真要保護 AMP 需改貼實體溫度感測器直讀 TAS5756M，勿再拿脫鉤的 SoC 當代理。
[Service]
Environment=MARVIN_TEMP_STOP_C=60
