# ops/com.antigravity.marvin.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  每日任務看門狗。每天 13:40 跑（在 daily 任務 12:00-13:00 之後）。
  純 stdlib → 故意用系統 /usr/bin/python3，避開 venv_simon 的 launchd dyld 死結。
  啟用：
    1. 把 <WEBHOOK> 換成 Discord webhook（要告警的話；不填就只寫 log）
    2. cp ops/com.antigravity.marvin.watchdog.plist ~/Library/LaunchAgents/
    3. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.antigravity.marvin.watchdog.plist
  停用：launchctl bootout gui/$(id -u)/com.antigravity.marvin.watchdog
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.marvin.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jackhuang/Code/Discord-voice-bot/scripts/cron_watchdog.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MARVIN_WATCHDOG_WEBHOOK</key>
        <string><WEBHOOK></string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>13</integer>
        <key>Minute</key><integer>40</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jackhuang/Library/Logs/Marvin/watchdog_cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jackhuang/Library/Logs/Marvin/watchdog_cron.log</string>
</dict>
</plist>
